# infrastructure/ansible/roles/hb-worker/templates/investmentology-price-refresh.timer.j2
[Unit]
Description=Investmentology Intraday Price Refresh Timer

[Timer]
OnCalendar={{ timer_price_refresh }}
Persistent=true

[Install]
WantedBy=timers.target
