# infrastructure/ansible/roles/hb-worker/templates/investmentology-post-screen.timer.j2
[Unit]
Description=Investmentology Post-Screen Analysis Timer

[Timer]
OnCalendar={{ timer_post_screen }}
Persistent=true

[Install]
WantedBy=timers.target
