# infrastructure/ansible/roles/hb-worker/templates/investmentology-premarket.timer.j2
[Unit]
Description=Investmentology Pre-Market Stop-Loss Timer

[Timer]
OnCalendar={{ timer_premarket }}
Persistent=true

[Install]
WantedBy=timers.target
