# infrastructure/ansible/roles/hb-worker/templates/investmentology-screen.timer.j2
[Unit]
Description=Investmentology Weekly Screen Timer

[Timer]
OnCalendar={{ timer_screen }}
Persistent=true

[Install]
WantedBy=timers.target
