# infrastructure/ansible/roles/hb-worker/templates/investmentology-daily-analyze.timer.j2
[Unit]
Description=Investmentology Daily Watchlist Analysis Timer

[Timer]
OnCalendar={{ timer_daily_analyze }}
Persistent=true

[Install]
WantedBy=timers.target
